# lxc/loadbalancer/rootfs/etc/init.d/start_proxy
#!/sbin/openrc-run

command="haproxy"
command_args="-D"
pidfile="/run/$RC_SVCNAME.pid"
command_background="yes"
