# webapp/start_webservice
#!/sbin/openrc-run

command="python3"
command_args="/home/api.py"
pidfile="/run/$RC_SVCNAME.pid"
command_background="yes"
depend()
{
	need localmount
}